# agentic_conversation_agent/rootfs/etc/services.d/agent/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting Agentic Conversation Agent server"

exec python3 -m agentic_conversation_agent.server
